# system/procoffee@.service
[Unit]
Description=Rodar processamento de audio

[Service]
ExecStart=/usr/local/bin/procoffee.sh %i

[Install]
WantedBy=multi-user.target

